# docker/cron/etc/crontabs/www-data
SHELL=/bin/sh
* * * * * test -f /local/www/artisan && /usr/local/bin/php /local/www/artisan schedule:run >> /dev/null 2>&1
# empty line